# etc/util/crontab/clean_twits_sh
#!/bin/bash

curl http://smeetv.com/etc/util/clean_twits.php > dev/null
