# etc/util/crontab/fetch_twits.sh
#!/bin/bash

curl http://smeetv.com/etc/util/fetch_twits.php > dev/null
