# etc/util/crontab/generate_stats.sh
#!/bin/bash

curl http://smeetv.com/status/generate/ > /dev/null
curl http://smeetv.com/status/stat/generate/ > /dev/null
curl http://smeetv.com/etc/util/index.stat.php > dev/null
