# etc/util/crontab/generate_featurehashtags.sh
#!/bin/bash

curl http://smeetv.com/etc/featuredhashtags/generate/index.php > /dev/null

