# etc/util/crontab/landing_page_generate.sh
#!/bin/bash

curl http://smeetv.com/etc/util/index.stat.php > dev/null
