# etc/util/crontab/landing_page_generate.sh
#!/bin/bash

curl http://smeetv.com/etc/util/landing_page_generatate_twits.php > /dev/null
