# etc/util/crontab/push_index_dump_priority.sh
#!/bin/bash

sleep 2;  curl http://smeetv.com/etc/util/index.keyword.php?priority=1 > /dev/null
sleep 10;  curl http://smeetv.com/etc/util/index.keyword.php?priority=2 > /dev/null
sleep 20;  curl http://smeetv.com/etc/util/index.keyword.php?priority=3 > /dev/null
sleep 30;  curl http://smeetv.com/etc/util/index.keyword.php?priority=4 > /dev/null
sleep 40;  curl http://smeetv.com/etc/util/index.keyword.php?priority=5 > /dev/null
sleep 50;  curl http://smeetv.com/etc/util/index.keyword.php?priority=6 > /dev/null
sleep 60;  curl http://smeetv.com/etc/util/index.keyword.php?priority=7 > /dev/null
sleep 70;  curl http://smeetv.com/etc/util/index.keyword.php?priority=8 > /dev/null
sleep 80;  curl http://smeetv.com/etc/util/index.keyword.php?priority=9 > /dev/null
sleep 90;  curl http://smeetv.com/etc/util/index.keyword.php?priority=10 > /dev/null


