# etc/util/crontab/push_index_dump_normal.sh
#!/bin/bash

curl http://fs.f1vlad.com/smeetv/smeetv.php/index.onequery.php?priority=0 > /dev/null 
curl http://fs.f1vlad.com/smeetv/smeetv.php/index.onequery.php?priority=1 > /dev/null 
curl http://fs.f1vlad.com/smeetv/smeetv.php/index.onequery.php?priority=2 > /dev/null



