# etc/util/crontab/push_index_dump_priority.sh
#!/bin/bash

sleep 2;  curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=1 > /dev/null
sleep 10; curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=2 > /dev/null
sleep 20; curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=3 > /dev/null
sleep 40; curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=4 > /dev/null
sleep 50; curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=5 > /dev/null 
sleep 70; curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=6 > /dev/null
sleep 80; curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=7 > /dev/null
sleep 100;curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=8 > /dev/null
sleep 110;curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=9 > /dev/null
sleep 113;curl http://fs.f1vlad.com/smeetv/smeetv.php/index.keyword.php?priority=10 > /dev/null


