# etc/util/crontab/push_index_dump_normal.sh
#!/bin/bash


curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null; 
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;

sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;

sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;


sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;

sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;

sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;

sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;





sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;



sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;



sleep 5;
curl http://smeetv.com/etc/util/index.onequery.php?priority=0 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=1 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=2 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=3 > /dev/null;
curl http://smeetv.com/etc/util/index.onequery.php?priority=4 > /dev/null;









