# etc/util/crontab/README


*/10 * * * * /var/www/smeetv.com/etc/util/crontab/landing_page_generate.sh > /dev/null
@daily /var/www/smeetv.com/etc/util/crontab/clean_twits.sh > /dev/null
@hourly /var/www/smeetv.com/etc/util/crontab/generate_featurehashtags.sh > /dev/null
*/30 * * * * /var/www/smeetv.com/etc/util/crontab/generate_stats.sh > /dev/null
*/30 * * * * /var/www/smeetv.com/etc/util/crontab/generate_stats_30_min_index.sh > /dev/null
* * * * * /var/www/smeetv.com/etc/util/crontab/push_index_dump_normal.sh > /dev/null
*/10 * * * * /var/www/smeetv.com/etc/util/crontab/push_index_dump_priority.sh  > /dev/null
*/15 * * * * /var/www/smeetv.com/etc/util/crontab/fetch_twits.sh > /dev/null



